# runner/build/linux/runner.service
[Unit]
Description=AgentMesh Runner - AI Agent Task Executor
Documentation=https://agentmesh.dev/docs/runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agentmesh
Group=agentmesh
ExecStart=/usr/local/bin/runner run --config /etc/agentmesh/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/lib/agentmesh /var/log/agentmesh

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentmesh-runner

[Install]
WantedBy=multi-user.target
